# roles/emacs/files/systemd/notes-sync.service
[Unit]
Description=Sync ~/Notes with Dropbox

[Service]
ExecStart=sh -c "$HOME/.local/bin/notes-sync"
