# roles/emacs/files/systemd/notes-sync.service
[Unit]
Description=Sync ~/Notes with Dropbox

[Service]
ExecStart=rclone sync ./Notes dropbox:Notes
