# roles/emacs/files/systemd/emacs.service.d/override.conf
[Service]
ExecReload=/usr/bin/emacsclient --eval "(load-file user-init-file)"