# roles/program.emacs/files/systemd-user/notes-sync.service
[Unit]
Description=Sync ~/Notes with Dropbox

[Service]
ExecStart=sh -c '~/.local/bin/notes-sync'
