# roles/program.emacs/files/systemd-user/notes-sync.service
[Unit]
Description=Sync ~/Notes.

[Service]
Environment=SSH_AUTH_SOCK=%t/keyring/ssh
ExecStart=sh -c '~/.local/bin/notes-sync'
