# ansible/roles/script.notes/files/systemd-user/notes-sync.timer
[Unit]
Description=Sync ~/Notes every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=default.target
