# roles/wm.eww/files/eww.service
[Unit]
Description=ElKowars wacky widgets

[Service]
ExecStart=sh -c "~/.local/bin/eww daemon --no-daemonize"

[Install]
WantedBy=default.target